# picam-control.service
[Unit]
Description=PiCam Control Agent (Tailscale-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ac
WorkingDirectory=/home/ac/ac-picam-control
Environment=PICAM_AGENT_HOST=127.0.0.1
Environment=PICAM_AGENT_PORT=8080
Environment=PICAM_SERVICE_NAME=picam
ExecStart=/usr/bin/python3 /home/ac/ac-picam-control/agent.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
